# templates/etc/profile.d/ansible.sh.j2
# {{ ansible_managed }}
export ANSIBLE_FORCE_COLOR="1"
export PY_COLORS="1"
